# scripts/tso-logging.service
[Unit]
Description=TSO Background Logging Daemon
Documentation=https://github.com/chukfinley/tso
After=network.target mysql.service mariadb.service apache2.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/serveros
Environment="TSO_INSTALL_DIR=/opt/serveros"
ExecStart=/usr/bin/php /opt/serveros/scripts/logging-daemon.php
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tso-logging

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/serveros/logs

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target

